# services/remote-desktop-gateway/systemd/systemd.service
[Unit]
Description=Remote Desktop Gateway
