# services/monitoring/systemd/systemd.service
[Unit]
Description=Monitoring (Grafrana, Prometheus)
