# services/dvr/systemd/systemd.service
[Unit]
Description=DVR Services (Torrents, Arr)
