# services/auth/systemd/systemd.service
[Unit]
Description=Auth Services (SSO, HTTP Basic Auth)
