# services/documentation/systemd/systemd.service
[Unit]
Description=Homelab Documentation Generator
